# apns-server/deploy/com.cccompanion.apns-server.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.cccompanion.apns-server</string>

  <key>ProgramArguments</key>
  <array>
    <string>/path/to/CcCompanion/apns-server/.venv/bin/python3</string>
    <string>/path/to/CcCompanion/apns-server/push.py</string>
    <string>--config</string>
    <string>/path/to/CcCompanion/apns-server/config.toml</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/path/to/CcCompanion/apns-server</string>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <dict>
    <key>SuccessfulExit</key>
    <false/>
    <key>Crashed</key>
    <true/>
  </dict>

  <key>ThrottleInterval</key>
  <integer>10</integer>

  <key>StandardOutPath</key>
  <string>/path/to/CcCompanion/apns-server/server.log</string>

  <key>StandardErrorPath</key>
  <string>/path/to/CcCompanion/apns-server/server.err.log</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    <key>PYTHONUNBUFFERED</key>
    <string>1</string>
  </dict>
</dict>
</plist>
